# deploy/vps/lms.service.d/socket.conf
# tsk-403: zero-downtime рестарт LMS через systemd socket activation.
# systemd (lms.socket) держит слушающий сокет 127.0.0.1:8000 через рестарты
# сервиса — во время `systemctl restart lms` соединения встают в очередь,
# а не получают connect refused (иначе nginx отдаёт клиенту 502).
#
# Базовый unit lms.service НЕ изменён; этот drop-in переопределяет ExecStart,
# чтобы uvicorn слушал унаследованный от systemd fd, а не биндил порт сам.
# Откат: удалить этот drop-in и lms.socket, daemon-reload, restart lms.
[Unit]
Requires=lms.socket
After=lms.socket

[Service]
ExecStart=
ExecStart=/opt/lms/venv/bin/uvicorn app.api.main:app --fd 3
