# deploy/vps/lms.service
[Unit]
Description=LMS Core API (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
User=app
Group=app
WorkingDirectory=/opt/lms
EnvironmentFile=/opt/lms/.env
ExecStart=/opt/lms/venv/bin/uvicorn app.api.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/lms/app.log
StandardError=append:/var/log/lms/app.log

[Install]
WantedBy=multi-user.target
